# systemd/aiops-qproxy-runner.service

[Unit]
Description=AIOps Q Proxy Runner (v2.4 merged-final)
After=network.target

[Service]
Type=simple
# 核心配置
Environment=Q_BIN=/home/ubuntu/.local/bin/q
Environment=QWORKDIR=/home/ubuntu/huixin/aiops/aiops-qproxy-v2.4
Environment=QCTX_DIR=/home/ubuntu/huixin/aiops/aiops-qproxy-v2.4/ctx/final
Environment=QLOG_DIR=/home/ubuntu/huixin/aiops/aiops-qproxy-v2.4/logs
Environment=QDATA_DIR=/home/ubuntu/huixin/aiops/aiops-qproxy-v2.4/data


# SOP 配置
Environment=Q_SOP_DIR=/home/ubuntu/huixin/aiops/aiops-qproxy-v2.4/ctx/sop
Environment=Q_SOP_PREPEND=1

# 输出控制
Environment=NO_COLOR=1
Environment=CLICOLOR=0
Environment=TERM=dumb

# Q CLI 自动信任设置
Environment=Q_MCP_AUTO_TRUST=true
Environment=Q_MCP_SKIP_TRUST_PROMPTS=true
Environment=Q_TOOLS_AUTO_TRUST=true

WorkingDirectory=/home/ubuntu/huixin/aiops/aiops-qproxy-v2.4
ExecStart=/home/ubuntu/huixin/aiops/aiops-qproxy-v2.4/bin/qproxy-runner --listen=:8080
Restart=always
RestartSec=2
User=ubuntu
Group=ubuntu

[Install]
WantedBy=multi-user.target
